# 99-ugreen-leds.rules
SUBSYSTEM=="block", KERNEL=="sd[a-z]", ACTION=="add|remove", RUN+="/usr/bin/systemctl --no-block start ugreen-leds.service"
